# supervisor.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=5
numprocs=1
startsecs=0
redirect_stderr=true
stdout_logfile=/var/log/supervisor/nginx.log

[program:backend]
command=gunicorn --bind 0.0.0.0:8000 app.main:app
directory=/app
autostart=true
autorestart=true
startretries=5
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/backend.log 